# conf/etc/s6/services/02-nginx/run
#!/bin/sh
service-handler "NGINX" "nginx -g 'daemon off;'"