# conf/etc/s6/services/01-phpfpm/run
#!/bin/sh
service-handler "PHP-FPM" "php-fpm7 -F"
